# deploy/dreamstudio.service
[Unit]
Description=DreamStudio - AI API Gateway Platform
Documentation=https://github.com/CrazyCherrys/DreamStudio
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=dreamstudio
Group=dreamstudio
WorkingDirectory=/opt/dreamstudio
ExecStart=/opt/dreamstudio/dreamstudio
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dreamstudio

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/dreamstudio

# Environment - Server configuration
# Modify these values to change listen address and port
Environment=GIN_MODE=release
Environment=SERVER_HOST=0.0.0.0
Environment=SERVER_PORT=8080

[Install]
WantedBy=multi-user.target
